# config/launchd/com.cainiao.inventory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cainiao.inventory</string>
  <key>WorkingDirectory</key>
  <string>/Users/chenyuecai/Documents/菜鸟云仓自动化操作</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/chenyuecai/.local/node/bin/node</string>
    <string>core/sync-cainiao-inventory.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/chenyuecai/.local/node/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/chenyuecai/Documents/菜鸟云仓自动化操作/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chenyuecai/Documents/菜鸟云仓自动化操作/logs/launchd.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
